# server_bot/bot.conf
[program:telegrambot]
directory = /home/ubuntu/aiogram-bot-template-master
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu", USER="ubuntu"

stderr_loglife=/home/ubuntu/aiogram-bot-template-master/loglife_err.log
stdout_loglife=/home/ubuntu/aiogram-bot-template-master/lexitoglife.log